# deploy/systemd/ensembleworks-term.service.d/10-memory.conf
# OOM containment / cgroup membership for the term service. Kept as a drop-in so
# all of the memory policy lives in one place; see ensembleworks.slice. term
# carries no MemoryLow and no cap — it gets "the rest" of the slice, so a runaway
# tmux process is throttled (slice MemoryHigh) then OOM-killed in-cgroup (slice
# MemoryMax) while the protected core survives.
# deploy/bootstrap-debian.sh writes this same file.
[Service]
Slice=ensembleworks.slice
MemoryAccounting=yes
