# deploy/systemd/ensembleworks-sync.service.d/10-memory.conf
# OOM containment / cgroup membership for the sync service. Kept as a drop-in so
# all of the memory policy lives in one place; see ensembleworks.slice. sync is a
# protected core service: MemoryLow reserves a working set the slice
# reclaims/OOM-kills LAST, so pressure lands on term (the elastic "rest").
# deploy/bootstrap-debian.sh writes this same file.
[Service]
Slice=ensembleworks.slice
MemoryAccounting=yes
MemoryLow=512M
